# systemd/autologin@.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin %I --noclear %%I $TERM
